# 2_5_3/app/services/ckan.script.conf
[program:ckan_configure]
command=/opt/ckan/configure_ckan
autostart=true
autorestart=false
